# argonONE.service
[Unit]
Description=ArgonONE

[Service]
ExecStart=/storage/.config/argonONE

[Install]
WantedBy=multi-user.target
